# includes/etc/cron.d/glpi-entities-sync
#
# Default GLPI to Netbox Entity Sync Job
#

*/15 * * * * ansible-playbook nofusscomputing.netbox.entities --extra-vars "@/root/vars.yaml" -v
