# includes/etc/cron.d/glpi-entities-sync
#
# Default GLPI to Netbox Entity Sync Job
#

*/15 * * * * ansible-playbook nofusscomputing.phpipam_scan_agent.agent --extra-vars "/root/vars.yaml" -v
